# tradingbot.service
[Unit]
Description=Trading Bot
After=network.target

[Service]
Type=simple
User=your_user
WorkingDirectory=/home/your_user/trading_bot
ExecStart=/usr/bin/python3 main.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target